# luci/applications/luci-app-slaac-dns-sync/root/etc/hotplug.d/iface/95-slaac-dns-sync
#!/bin/sh

case "$ACTION" in
	ifup|ifupdate|ifdown)
		/etc/init.d/slaac-dns-sync enabled 2>/dev/null || exit 0
		/etc/init.d/slaac-dns-sync reload >/dev/null 2>&1
		;;
esac
