# luci/applications/luci-app-slaac-dns-sync/root/etc/uci-defaults/95_luci-slaac-dns-sync
#!/bin/sh

mkdir -p /var/lock
/etc/init.d/slaac-dns-sync enable >/dev/null 2>&1 || true
rm -f /tmp/luci-indexcache
exit 0
